# home/.local/share/systemd/user/internet-metered.target
[Unit]
Description=current connection to the internet is metered
Conflicts=internet-unmetered.target