# home/.local/share/systemd/user/internet-unmetered.target
[Unit]
Description=current connection to the internet is not metered
Conflicts=internet-metered.target