# home/.config/systemd/user/unmetered-internet.target
[Unit]
Description=current connection to the internet is not metered